# autostart/autostart.desktop
[Desktop Entry]
Type=Application
Name=Speech Stress Detector
Exec=/home/pi/.config/autostart/SpeechStressDetector
